# conf/supervisor.conf
[program:uwsgi-app]
command = /usr/local/bin/uwsgi --ini /portal/conf/uwsgi.ini
